# homeassistant-addon/rootfs/etc/services.d/ha-mcp-plus/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting ha-mcp-plus..."

export HA_URL="http://supervisor/core"
export HA_TOKEN="${SUPERVISOR_TOKEN}"
export CONFIG_PATH="/config"
export MCP_PORT="9584"
export MCP_SECRET_PATH="$(bashio::config 'mcp_secret_path')"

# InfluxDB
export INFLUX_TOKEN="$(bashio::config 'influx_token')"
export INFLUX_ORG="$(bashio::config 'influx_org')"
export INFLUX_BUCKET="$(bashio::config 'influx_bucket')"

# Grafana
export GRAFANA_TOKEN="$(bashio::config 'grafana_token')"

# Node-RED
export NODERED_TOKEN="$(bashio::config 'nodered_token')"

bashio::log.info "MCP endpoint: 0.0.0.0:9584${MCP_SECRET_PATH}"
bashio::log.info "Auto-discovering installed addons..."

cd /app
exec python -u server.py
